# systemd/fs-update_portfolio.service
[Unit]
Description=FundScope: T212 portfolio sync + git push
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/fundscope
EnvironmentFile=/home/ubuntu/fundscope/.env
ExecStart=/usr/local/bin/fs-run ingest/update_portfolio.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fs-update_portfolio
MemoryMax=600M
CPUQuota=80%
TimeoutStartSec=30min
