# systemd/fs-update_markets.service
[Unit]
Description=FundScope: refresh sector markets
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/fundscope
EnvironmentFile=/home/ubuntu/fundscope/.env
ExecStartPre=/bin/test -x /home/ubuntu/fundscope/.venv/bin/python
ExecStart=/usr/local/bin/fs-run ingest/update_markets.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fs-update_markets
MemoryMax=400M
CPUQuota=60%
TimeoutStartSec=10min
